# docker/vita.06-DEC-2016/_etc_supervisor.conf
[supervisord]
nodaemon=true

[program:mongo]
command=mongod

[program:sl-pm]
command=sl-pm --skip-default-install

;[include]
;files = /etc/supervisor/conf.d/*.conf
